# etc/profile.d/google_drive_ocaml_sync.sh
google-drive-ocamlfuse
google-drive-ocamlfuse -label personal ~/Desktop/Google_Drive
google-drive-ocamlfuse -label utexas ~/Desktop/UTexas_Drive
